# systemd/autotrader-saida-v2.service
[Unit]
Description=AUTOTRADER-SAIDA V2 - Worker (5min)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/roteiro_ds/AUTOTRADER-SAIDA-V2
Environment=DATA_DIR=/home/roteiro_ds/AUTOTRADER-SAIDA-V2/data
Environment=GAIN_MIN=3
ExecStart=/usr/bin/python3 /home/roteiro_ds/AUTOTRADER-SAIDA-V2/worker_saida_v2.py

[Install]
WantedBy=multi-user.target
